# deploy/systemd/user/sora-secretary-api.service
[Unit]
Description=SORA Secretary UI API
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/okota/code/sora-secretary
ExecStart=/home/okota/code/sora-secretary/.venv/bin/uvicorn api_server:app --host 127.0.0.1 --port 8787
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
